# coolify-mcp-server.service
[Unit]
Description=Coolify MCP Server
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/coolify-mcp-server
Environment="PATH=/home/YOUR_USERNAME/coolify-mcp-server/venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/YOUR_USERNAME/coolify-mcp-server/venv/bin/python server.py
Restart=always
RestartSec=10

# Load environment variables from .env
EnvironmentFile=/home/YOUR_USERNAME/coolify-mcp-server/.env

[Install]
WantedBy=multi-user.target
